# com.willow.toss-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.willow.toss-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/dongwookkim/scripts/drive-launcher.sh</string>
        <string>--cwd</string>
        <string>/Volumes/PRO-G40/app-dev/willow-invt</string>
        <string>/Volumes/PRO-G40/app-dev/willow-invt/scripts/run-toss-sync.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/dongwookkim</string>

    <!-- 매일 KST 16:10 (한국장 마감 후) + 06:40 (미국장 마감 후) -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>16</integer>
            <key>Minute</key>
            <integer>10</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>6</integer>
            <key>Minute</key>
            <integer>40</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/dongwookkim/logs/toss-sync/launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/dongwookkim/logs/toss-sync/launchd.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/dongwookkim</string>
    </dict>
</dict>
</plist>
