# com.willow.gmail-auto-label.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.willow.gmail-auto-label</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Volumes/PRO-G40/app-dev/willow-invt/scripts/run-gmail-auto-label.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Volumes/PRO-G40/app-dev/willow-invt</string>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>StandardOutPath</key>
    <string>/Volumes/PRO-G40/app-dev/willow-invt/scripts/logs/gmail-auto-label.log</string>
    <key>StandardErrorPath</key>
    <string>/Volumes/PRO-G40/app-dev/willow-invt/scripts/logs/gmail-auto-label.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
